# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

# --- glottos-courses : Next.js standalone (node server.js) ---
[program:courses]
directory=/app/courses
command=node server.js
environment=NODE_ENV="production",PORT="8080",HOSTNAME="127.0.0.1",DATABASE_URL="%(ENV_DATABASE_URL_COURSES)s"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- text-tutor : Express via tsx ---
[program:tutor]
directory=/app/tutor
command=npm start
environment=NODE_ENV="production",PORT="4000",DATABASE_URL="%(ENV_DATABASE_URL_TUTOR)s"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- nginx edge router ---
[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
